# netwatchdog.service
[Unit]
Description=NetWatchDog v5.1 - Self Healing WiFi and Service Watchdog
After=network-online.target usb-wifi-wpa.service zerotier-one.service mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netwatchdog
ExecStart=/usr/bin/python3 /opt/netwatchdog/netwatchdog_standby.py
Restart=always
RestartSec=10
MemoryMax=128M

[Install]
WantedBy=multi-user.target
